# supervisor-services.ini
[program:portal]
directory=/app/portal
command=node server.js
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/portal.log
stderr_logfile=/var/log/supervisor/portal_error.log
environment=PORT=8080,NODE_ENV=production

[program:parenting]
directory=/app/parenting
command=node server.js
autostart=true
autorestart=true
priority=2
stdout_logfile=/var/log/supervisor/parenting.log
stderr_logfile=/var/log/supervisor/parenting_error.log
environment=PORT=3000,NODE_ENV=production

[program:cello]
directory=/app/cello-practise
command=node server.js
autostart=true
autorestart=true
priority=3
stdout_logfile=/var/log/supervisor/cello.log
stderr_logfile=/var/log/supervisor/cello_error.log
environment=PORT=3001,NODE_ENV=production
